# .config/systemd/user/xdg-desktop-portal.service.d/override.conf
# https://wiki.archlinux.org/title/XDG_Desktop_Portal
# https://wiki.archlinux.org/title/Xdg-utils
# suggest gnome despite using standalone wm

[Service]
Environment="XDG_CURRENT_DESKTOP=GNOME"
