# .hardware/99-cadmouse.rules
KERNEL=="hidraw*", ATTRS{idVendor}=="256f", ATTRS{idProduct}=="c650", RUN+="/usr/local/bin/cadmousectl -s 50 -S 0"
